# config/com.anki-renderer.deploy-webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anki-renderer.deploy-webhook</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ben/services/anki-renderer/deploy-webhook.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DEPLOY_SECRET</key>
        <string>__DEPLOY_SECRET__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ben/services/data/anki-renderer-deploy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ben/services/data/anki-renderer-deploy.log</string>
</dict>
</plist>
